# voice/systemd/serena.service
[Unit]
Description=Serena — Always-On Voice AI Assistant
After=graphical-session.target
Wants=pulseaudio.service pipewire.service

[Service]
Type=simple
ExecStart=/home/raghav/Documents/Projects/serena/serena/.venv/bin/python -m serena.main
WorkingDirectory=/home/raghav/Documents/Projects/serena/serena
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Audio access
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native

[Install]
WantedBy=default.target
